# rtorrent/root/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

cd /app/flood || exit

exec \
	s6-setuidgid abc /usr/bin/npm start
